# supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=python3 manage.py runserver 0.0.0.0:8000

[program:celery]
command=celery -A mysite worker -l INFO